# root/etc/services.d/nginx/run
#!/bin/sh
# Papra nginx s6 service.
# Serves the built SolidStart/Vite SPA on port 1221 and reverse-proxies /api/
# to the node backend on 127.0.0.1:1222.

echo "[INFO] Starting nginx..."
exec nginx -g "daemon off;" -c /usr/local/etc/nginx/papra.conf
